# linux/hvac-bridge.service
[Unit]
Description=HVAC Bridge Daemon (RS485 polling, SensorFlags, MQTT, Bridge RPC)
After=network-online.target
Wants=network-online.target

[Service]
# Type=notify so we can use the systemd watchdog: bridge_daemon calls
# sd_notify("READY=1") once it enters its polling loop and sd_notify(
# "WATCHDOG=1") on every iteration.  If the loop ever blocks (deadlock,
# RS485 lockup, MQTT thread wedged) for longer than WatchdogSec, systemd
# kills + restarts us via the Restart=always policy below.
Type=notify
WatchdogSec=60
NotifyAccess=main
User=arduino
WorkingDirectory=/home/arduino/hvac-controller/linux

# Full path to the daemon script — adjust if repo is cloned elsewhere
ExecStart=/usr/bin/python3 /home/arduino/hvac-controller/linux/bridge_daemon.py

# Restart policy: always restart, 10-second delay, giving MQTT/RS485 time to recover
Restart=always
RestartSec=10

# Standard output and error go to the systemd journal
# View with: journalctl -u hvac-bridge -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hvac-bridge

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
